# rootfs/etc/systemd/system/garnet-rtc-save.timer
[Unit]
Description=garnet: periodically save RTC offset

[Timer]
OnBootSec=10min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
